# root/etc/services.d/unbound/run
#!/usr/bin/with-contenv bash
cd /config || exit
exec \
	/usr/sbin/unbound -d -v -c /etc/unbound/unbound.conf